# media-volume-resize.service
[Unit]
After=local-fs.target

[Service]
ExecStart=/sbin/resize2fs /dev/xvdd

[Install]
WantedBy=default.target